# systemd/blueman-applet.service
[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=/usr/bin/env blueman-applet
Restart=always
RestartSec=5

[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Description=Blueman-Applet: A tray applet for managing bluetooth
Documentation=man:blueman-applet(1)
